# ros2/robot-auto.service
[Unit]
Description=Robot Auto System with LiDAR and SLAM
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/sbx/ros2
ExecStart=/sbx/ros2/start_robot.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=ROS_DOMAIN_ID=0
Environment=DISPLAY=:0

# Give the service time to start up
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target 